# debian/go-apt-cache.service
[Unit]
Description=Go APT Cache (RFC-compliant caching proxy)
Documentation=https://github.com/yolkispalkis/go-apt-cache
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/go-apt-cache -config /etc/go-apt-cache/config.yaml
User=go-apt-cache
Group=go-apt-cache
WorkingDirectory=/var/cache/go-apt-cache
RuntimeDirectory=go-apt-cache
RuntimeDirectoryMode=0750

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true

# Restart policy
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target